# 3rdparty/arc-service
#!/bin/bash

### BEGIN INIT INFO
# Provides: arc-service-@@name@@
# Required-Start: $local_fs $syslog $remote_fs
# Required-Stop: $local_fs $syslog $remote_fs
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: Daemon for @@name@@
# Description: Connect to an android device
### END INIT INFO

case "$1" in
    start)
	          echo "Starting AndroidRemoteControl Daemon for @@name@@..."
            adb connect @@address@@
            echo $! > /run/ar-@@name@@.pid
            ;;
    stop)
            echo "Stopping AndroidRemoteControl Daemon for @@name@@..."
            adb disconnect @@address@@
            rm /run/ar-@@name@@.pid
            ;;
    restart)
            $0 stop
            sleep 1
            $0 start
            ;;
    kill)
            #pkill -f speaker.sh
            kill -9 `cat /run/ar-@@name@@.pid`
            ;;
    *)
            echo "Usage: $0 {start|stop|restart|kill}"
            exit 1
            ;;
esac
